# src/crontab
0       1       *       *       *       python /app/deletion.py
0       2       *       *       *       python /app/concat.py
